# deploy/gongxin-monitor.service
[Unit]
Description=Gongxin platform status monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/gongxin-monitor
EnvironmentFile=/opt/gongxin-monitor/.env
ExecStart=/usr/bin/env node /opt/gongxin-monitor/src/monitor.js
User=gongxin-monitor
Group=gongxin-monitor
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
